# home/Library/LaunchAgents/local.downloads-to-desktop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.downloads-to-desktop</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-l</string>
    <string>-c</string>
    <!-- fswatch を Homebrew PATH で解決するためログインシェルで起動 -->
    <string>exec "$HOME/.local/bin/downloads-to-desktop.sh"</string>
  </array>

  <!-- ログイン時に自動起動し、落ちたら再起動 -->
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/downloads-to-desktop.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/downloads-to-desktop.log</string>
</dict>
</plist>
