# home/Library/LaunchAgents/local.downloads-to-desktop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.downloads-to-desktop</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-l</string>
    <string>-c</string>
    <!-- osascript 経由で Finder に move を依頼。TCC は Finder が解決するため bash/osascript への FDA 付与は不要 -->
    <string>exec /usr/bin/osascript "$HOME/.local/lib/downloads-to-desktop.applescript"</string>
  </array>

  <!-- ログイン時に1回 + 5秒ごとに polling -->
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>5</integer>

  <key>StandardOutPath</key>
  <string>/tmp/downloads-to-desktop.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/downloads-to-desktop.log</string>
</dict>
</plist>
