# home/Library/LaunchAgents/local.brew-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.brew-update</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-l</string>
    <string>-c</string>
    <string>exec "$HOME/dotfiles/scripts/darwin/brew_update.sh"</string>
  </array>

  <!-- 毎日 10:00 に実行。スリープ中は StartCalendarInterval が発火せず次の wake に持ち越され、
       45-180 秒の DarkWake (PowerNap) 中に走ると大型 Cask の DL が Back to Sleep で切れて
       make homebrew が失敗する。フル覚醒している日中に寄せるため 07:00 から移動 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/brew-update.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/brew-update.log</string>
</dict>
</plist>
